# ggzd/debian/ggzd.postinst
#!/bin/sh

set -e

case "$1" in
	configure*)
	#ldconfig
	update-rc.d ggzd defaults >/dev/null
	mkdir /var/ggzd || true 2>/dev/null
	chown -R games.games /var/ggzd
	chmod 755 /var/ggzd
	/etc/init.d/ggzd start || true
	;;
	*)
	;;
esac

#DEBHELPER#

exit 0
