# ggzd/debian/ggzd.postrm
#!/bin/sh

set -e

if [ purge = "$1" ]; then
update-rc.d ggzd remove >/dev/null
rm -rf /var/ggzd
fi

#DEBHELPER#

exit 0
